# com.powerbhoomi.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Identification -->
    <key>Label</key>
    <string>com.powerbhoomi.app</string>
    
    <!-- Program to Run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/Applications/POWER-BHOOMI/app/venv/bin/activate && python3 ~/Applications/POWER-BHOOMI/app/bhoomi_web_app_v2.py</string>
    </array>
    
    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/USER/Applications/POWER-BHOOMI/app</string>
    
    <!-- Run at Load -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep Alive (restart if crashes) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <!-- Restart Configuration -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <!-- Standard Output -->
    <key>StandardOutPath</key>
    <string>/Users/USER/Library/Logs/POWER-BHOOMI/stdout.log</string>
    
    <!-- Standard Error -->
    <key>StandardErrorPath</key>
    <string>/Users/USER/Library/Logs/POWER-BHOOMI/stderr.log</string>
    
    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>DATA_DIR</key>
        <string>/Users/USER/Documents/POWER-BHOOMI</string>
        <key>LOG_DIR</key>
        <string>/Users/USER/Library/Logs/POWER-BHOOMI</string>
        <key>CONFIG_DIR</key>
        <string>/Users/USER/.config/power-bhoomi</string>
    </dict>
    
    <!-- Process Type -->
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <!-- Network Access -->
    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SockServiceName</key>
            <string>5001</string>
            <key>SockType</key>
            <string>stream</string>
            <key>SockFamily</key>
            <string>IPv4</string>
        </dict>
    </dict>
    
    <!-- Resource Limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
        <key>NumberOfProcesses</key>
        <integer>100</integer>
    </dict>
    
    <!-- Nice Value (Lower = Higher Priority) -->
    <key>Nice</key>
    <integer>0</integer>
    
    <!-- Run only when user is logged in -->
    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
        <string>Background</string>
        <string>LoginWindow</string>
        <string>StandardIO</string>
    </array>
</dict>
</plist>

